# project/rfs/etc/rc.local
#!/bin/sh -e
/bin/bash /usr/local/etc/rc.local 2>&1 >/tmp/rc.local.log
exit 0

